# crontab
# Deutsche Bahn Verbindungsüberwachung - Docker Cron
# Läuft alle 6 Stunden: 6:00, 12:00, 18:00, 00:00

# Umgebungsvariablen für Cron
PATH=/usr/local/bin:/usr/bin:/bin
PYTHONPATH=/app

# Cron-Jobs (4x täglich)
0 6 * * * cd /app && python src/main.py --run >> /var/log/bahnabfrage/cron.log 2>&1
0 12 * * * cd /app && python src/main.py --run >> /var/log/bahnabfrage/cron.log 2>&1
0 18 * * * cd /app && python src/main.py --run >> /var/log/bahnabfrage/cron.log 2>&1
0 0 * * * cd /app && python src/main.py --run >> /var/log/bahnabfrage/cron.log 2>&1

# Leere Zeile am Ende (wichtig für Cron)