# crontab
# Deutsche Bahn Verbindungsüberwachung - Docker Cron  
# Läuft alle 3 Minuten für sofortige Erkennung neuer Verbindungen

# Umgebungsvariablen für Cron - Essential für Container
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
PYTHONPATH=/app

# Cron-Job (alle 3 Minuten) - Überwachung zukünftiger Verbindungen
*/3 * * * * bahnmonitor /bin/bash -c 'echo "=== $(date) - 3-Min Check gestartet ===" >> /var/log/bahnabfrage/cron.log; cd /app && source .env 2>/dev/null && python src/main.py --run >> /var/log/bahnabfrage/cron.log 2>&1'

# Leere Zeile am Ende (wichtig für Cron)